# ansible/roles/test_monitoring/files/test_monitor.service
[Unit]
Description=Monitor the "test" process
After=network.target

[Service]
Type=oneshot
ExecStart=/home/ubuntu/var/scripts/monitor_test.sh